# sys-apps/ps4fancontrol/files/ps4fancontrol.service
[Unit]
Description=Adjust threshold temperature

[Service]
ExecStart=/sbin/ps4fancontrol --no-gui

[Install]
WantedBy=multi-user.target
